# .config/systemd/user/home-robert-bin8.automount
[Unit]
Description=Automount via systemd of FREENASPassThru bin8

[Automount]
Where=/home/robert/bin8
DirectoryMode=0766
TimeoutIdleSec=7200

[Install]
WantedBy=multi-user.target
