# .config/systemd/system/home-robert-repos.automount
[Unit]
Description=Automount via systemd of ~/repos/ 

[Automount]
Where=/home/robert/repos
DirectoryMode=0766
TimeoutIdleSec=7200

[Install]
WantedBy=multi-user.target
